# systemd/zfs-autosnap-snap.timer
[Unit]
Description=Take ZFS snapshot hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
