# systemd/zfs-autosnap-gc.service
[Unit]
Description=Garbage collect ZFS auto snapshots

[Service]
ExecStart=zfs-autosnap gc
