# systemd/zfs-autosnap-snap.service
[Unit]
Description=Take ZFS auto snapshot

[Service]
ExecStart=zfs-autosnap snap
