# systemd/zfs-autosnap-gc.timer
[Unit]
Description=Garbage collect ZFS auto snapshots daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
